# src/ReferenceRAG.Service/scripts/Archived/obsidianrag.service
[Unit]
Description=ReferenceRAG Knowledge Base Service
Documentation=https://github.com/your-repo/ReferenceRAG
After=network.target

[Service]
Type=simple
# The executable path - change this to match your installation
ExecStart=/usr/bin/dotnet /opt/referencerag/ReferenceRAG.Service.dll
WorkingDirectory=/opt/referencerag

# User and group - service will NOT run as root
User=referencerag
Group=referencerag

# Restart policy
Restart=on-failure
RestartSec=10

# Environment variables
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_RUNNING_IN_CONTAINER=false

# Logging
StandardOutput=append:/var/log/referencerag/service.log
StandardError=append:/var/log/referencerag/error.log

# Graceful shutdown
TimeoutStopSec=30

# Hardening - don't allow core dumps
LimitCORE=0

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/referencerag/data /var/log/referencerag

[Install]
WantedBy=multi-user.target
